# json.queried.revs-daily.cron
#!/bin/bash
# by: lilong'en(lilongen@163.com)
#

cd "#CCV_HOME#"
perl structuralizeLogInfo2Json.pl
